# cron.daily/etc-autorepack
#!/bin/sh
#FIXME: support other directories than /etc/
cd /etc
git-repack -a -d
